# roles/workshop/files/workshop.service
[Unit]

[Service]
Restart=no
ExecStart=/usr/bin/python2 /tmp/bin/workshop.py

[Install]
WantedBy=multi-user.target
